# roles/shell/files/etc/cron.hourly/znc_saveconfig.sh
#!/bin/dash

PID=$(pgrep -u znc znc)

#Rewrite znc.conf
[ $PID ] && kill -s USR1 $PID
